# .config/systemd/user/dummy-mic.service
[Unit]
Description=Dummy Microphone (PipeWire)
After=pipewire.service

[Service]
Type=simple
ExecStart=pw-loopback --capture-props="media.class=Audio/Source node.name=dummy_mic node.description='Dummy Microphone'"
Restart=always
RestartSec=1

[Install]
WantedBy=default.target   
